# pkg/crawl.install
post_install() {
    cat <<'EOF'
==> crawl installed.
==> Add your user to the 'video' group for brightness control:
==>   sudo usermod -aG video $USER
==>
==> Enable the user service:
==>   systemctl --user enable --now crawl
EOF
}

post_upgrade() {
    post_install
}
